# rootfs/etc/s6-overlay/s6-rc.d/reload/run
#!/command/with-contenv sh

exec /reload.sh