# rootfs/etc/s6-overlay/s6-rc.d/login/up
#!/command/execlineb -P

with-contenv
exec /login.sh
