# rootfs/etc/s6-overlay/s6-rc.d/webdav/run
#!/command/with-contenv sh

exec /webdav.sh